# roles/container-host/files/rootfs/etc/systemd/system/snapd.service.d/override.conf
[Unit]
Requires=drbd.service
RequiresMountsFor=/var/snap /var/lib/snapd
