# roles/container-host/files/etc/systemd/system/snap.lxd.daemon.service.d/override.conf
[Unit]
RequiresMountsFor=/var/snap /var/lib/snapd
Requires=snapd
