# roles/container-host/files/rootfs/etc/systemd/system/docker.service.d/override.conf
[Unit]
Requires=drbd.service
RequiresMountsFor=/var/lib/docker
