# roles/nfs-server/files/etc/systemd/system/nfs-server.service.d/override.conf
[Unit]
Requires=drbd.service
RequiresMountsFor=/mnt/nfs
